# 00.Installation/template/systemd/deluged.service
[Unit]
Description=Deluge Daemon Service
After=network.target

[Service]
LimitNOFILE=666666
User=root
ExecStart=/usr/bin/deluged -d -l /etc/inexistence/01.Log/deluged.log -L info
ExecStop=/usr/bin/killall -w -s 9 /usr/bin/deluged
Restart=on-failure
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
